# zabbix-frontend/supervisor.conf
[supervisord]
nodaemon=true
childlogdir=/tmp/

[program:nginx]
command=/usr/sbin/nginx
autorestart=true

[program:php5-fpm]
command=/usr/sbin/php5-fpm --nodaemonize
autorestart=true
